# logvue/etc/supervisord.conf
#  Copyright 2018, 2019 Liang Chen <liangchenomc@gmail.com>
#
#  This file is part of kitkit.
#
#  kitkit is free software: you can redistribute it and/or modify it
#  under the terms of the GNU Affero General Public License
#  as published by the Free Software Foundation, either version 3 of
#  the License, or (at your option) any later version.
#
#  kitkit is distributed in the hope that it will be useful,
#  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  GNU Affero General Public License for more details.
#
#  You should have received a copy of the GNU Affero General Public License
#  along with kitkit. If not, see <https://www.gnu.org/licenses/>.

[supervisord]
nodaemon=true
pidfile=%(ENV_HOME)s/temp/supervisord.pid
logfile=%(ENV_HOME)s/temp/pm.log
identifier=gray
childlogdir=%(ENV_HOME)s/temp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=%(ENV_HOME)s/temp/supervisor.sock

[supervisorctl]
serverurl=unix://%(ENV_HOME)s/temp/supervisor.sock

[program:mongoserver]
command=mongod --dbpath %(ENV_HOME)s/data/mongo
priority=100

[program:elastic]
command=/usr/share/elasticsearch/bin/elasticsearch -Edefault.path.conf=%(ENV_HOME)s/etc -Epath.data=%(ENV_HOME)s/data/elastic -Epath.logs=%(ENV_HOME)s/data/temp
environment=ES_JAVA_OPTS="-Dlog4j2.disable.jmx=true"
priority=300

#[program:amqphub]
#command=rabbitmq-server
#priority=100

#[program:stashlog]
#command=logstash
#priority=350

[program:gray]
command=/usr/local/bin/graylog.sh
priority=500
